# chef_repo/site-cookbooks/application_nodejs/templates/default/nodejs.systemd.service.erb
[Service]
ExecStart=<%= @npm_dir %>/bin/npm start
Restart=always
WorkingDirectory=<%= @app_dir %>
StandardOutput=syslog
SyslogIdentifier=<%= @app_name %>
User=<%= @user %>
Group=<%= @group %>
<% @environment.each do |key, value| -%>
Environment=<%= key %>=<%= value %>
<% end -%>

[Install]
WantedBy=multi-user.target
